# deploy/order.service
[Unit]
Description=Order System (Gunicorn)
After=network.target postgresql.service

[Service]
User=www-data
Group=www-data
WorkingDirectory=/opt/order
EnvironmentFile=/opt/order/.env
ExecStart=/opt/order/venv/bin/gunicorn config.wsgi:application -c /opt/order/deploy/gunicorn.conf.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
